# roles/xray/templates/xray.service.j2
[Unit]
Description=Xray Service
Documentation=https://github.com/xtls/Xray-core
After=network.target nss-lookup.target

[Service]
User={{ xray_user }}
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ExecStart={{ xray_install_dir }}/xray run -config {{ xray_config_dir }}/config.json
Restart=on-failure
RestartPreventExitStatus=23
LimitNPROC=10000
LimitNOFILE=1000000

[Install]
WantedBy=multi-user.target
